# packages/hotspotos-core/files/etc/uci-defaults/99-hotspotos
#!/bin/sh
# HotspotOS UCI Defaults
# Initial setup on first boot

mkdir -p /etc/hotspotos
mkdir -p /var/log
mkdir -p /var/run
mkdir -p /www/hotspotos

chmod 755 /etc/init.d/hotspotos
chmod 755 /usr/lib/hotspotos/*.sh
chmod 755 /usr/sbin/hotspotos

/etc/init.d/hotspotos enable

uci -q batch <<-EOF >/dev/null
    set uhttpd.main.listen_http='0.0.0.0:8080'
    set uhttpd.main.listen_https='0.0.0.0:443'
    set uhttpd.main.home='/www'
    commit uhttpd
EOF

uci -q batch <<-EOF >/dev/null
    set network.lan.ipaddr='192.168.1.20'
    commit network
EOF

/etc/init.d/uhttpd restart

echo "HotspotOS initial setup completed"
